# .config/systemd/user/llama-server.service
[Unit]
Description=llama.cpp inference server
After=network-online.target
Wants=network-online.target
BindsTo=mcp-servers.service
After=mcp-servers.service

[Service]
Type=simple
Environment=MODEL=/home/romedius/Documents/models/gemma-4-26B-A4B-it-heretic.Q8_0.gguf
Environment=PORT=1337
Environment=CONTEXT=131072
Environment=GPU_LAYERS=999
EnvironmentFile=-%h/.config/llama-server/env

ExecStart=/usr/bin/distrobox-enter hyprome-dev-distrobox-quadlet -- \
  /home/linuxbrew/.linuxbrew/bin/llama-server \
  --model ${MODEL} \
  --host 0.0.0.0 \
  --port ${PORT} \
  --webui-mcp-proxy \
  --ctx-size ${CONTEXT} \
  --n-gpu-layers ${GPU_LAYERS}

Restart=on-failure
RestartSec=10
TimeoutStopSec=30

[Install]
WantedBy=default.target
