# .config/systemd/user/reset-hyprome.service
[Unit]
Description=Reset and restart hyprome container

[Service]
Type=oneshot
ExecStart=%h/.local/bin/reset-hyprome-dev
