# .config/systemd/user/reset-hyprome.service
[Unit]
Description=Reset hyprome dev container

[Service]
Type=oneshot
ExecStart=reset-hyprome-dev
RemainAfterExit=true
